# supervisor.conf
; =============================================================
; Supervisor — Caron Immobilier
; Copier dans /etc/supervisor/conf.d/caron.conf
; Puis : supervisorctl reread && supervisorctl update
; =============================================================

[program:caron-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/caron/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
directory=/var/www/caron
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/caron-worker.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=3600

[program:caron-scheduler]
; Lance le scheduler Laravel toutes les minutes (remplace la ligne crontab)
process_name=%(program_name)s
command=php /var/www/caron/artisan schedule:work
directory=/var/www/caron
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/caron-scheduler.log
stdout_logfile_maxbytes=5MB
